# deployment/systemd/policyd-rate-guard-cleanup.timer
[Unit]
Description=Timer for PolicydRateGuard cleanup job

[Timer]
OnCalendar=00:00:10
Persistent=true

[Install]
WantedBy=timers.target
