# deployment/systemd/policyd-rate-guard-cleanup.timer
[Unit]
Description=Timer for PolicydRateGuard cleanup job

[Timer]
OnCalendar=00:00:00

[Install]
WantedBy=timers.target
